# conf/task.cron
# Every 2 hours
0 */2 * * * root __INSTALL_DIR__/appstore/maintenance.sh
