# conf/task.cron
# Every 4 hours
0 */4 * * * root __INSTALL_DIR__/update_catalog.sh

# Every day
0 */4 * * * root __INSTALL_DIR__/update_catalog.sh
